# savehistory.service
[Unit]
Description=Save my firefox history from inevitable deletion

[Service]
Type=simple
ExecStart=%h/projects/savehistory/save.sh
